# crontab/crontab
# https://crontab.guru/

0 12 * * * php bin/console app:send-email-appoinment-reminder
